# systemd/postgresql.target
[Unit]
Description=PostgreSQL RDBMS

[Install]
WantedBy=multi-user.target
